# reliance_ams_local-master/ams-start.service
[Unit]
Description=Enable AMS apps

[Service]
ExecStart=/home/ams-core/ams-start.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
